# archive/telegram-bot.service
[Unit]
Description=Telegram Gemini Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bot
Environment="TG_BOT_TOKEN=YOUR_TOKEN_HERE"
Environment="GEMINI_API_KEY=YOUR_KEY_HERE"
ExecStart=/root/bot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
